# crontab
*/2 * * * * root /usr/bin/python3 /app/app.py >> /var/log/insta.log
